# supervisord.conf/flask_bing.conf
[program:bing]  # risk是自己给进程取的名，随意
command=gunicorn -c gunicorn.conf.py manage:flask_app
;directory=/data/flask_bing  # 项目目录
startsecs=0 
stopwaitsecs=0   
autostart=true
autorestart=true