# supervisord.conf/flask_bing.conf
[program:bing]
command=gunicorn -c gunicorn.conf.py manage:flask_app
;directory=/data/flask_bing  # 项目目录
startsecs=0 
stopwaitsecs=0   
autostart=true
autorestart=true